# test/data/cronie-1.4.10-7.fc19.x86_64.filelist
rw-r--r--,"","",2,root,root,128,/etc/cron.d/0hourly
rw-r--r--,"",cn,3,root,root,0,/etc/cron.deny
rw-r--r--,"",cn,4,root,root,293,/etc/pam.d/crond
rw-r--r--,"",cn,5,root,root,110,/etc/sysconfig/crond
rw-r--r--,"","",6,root,root,263,/lib/systemd/system/crond.service
rwsr-xr-x,"","",7,root,root,57568,/usr/bin/crontab
rwxr-xr-x,"","",8,root,root,66024,/usr/sbin/crond
rw-r--r--,"",d,10,root,root,48,/usr/share/doc/cronie-1.4.10/AUTHORS
rw-r--r--,"",d,11,root,root,4091,/usr/share/doc/cronie-1.4.10/COPYING
rw-r--r--,"",d,12,root,root,69610,/usr/share/doc/cronie-1.4.10/ChangeLog
rw-r--r--,"",d,13,root,root,364,/usr/share/doc/cronie-1.4.10/INSTALL
rw-r--r--,"",d,14,root,root,543,/usr/share/doc/cronie-1.4.10/README
rw-r--r--,"",d,15,root,root,2618,/usr/share/man/man1/crontab.1.gz
rw-r--r--,"",d,16,root,root,4980,/usr/share/man/man5/crontab.5.gz
rw-r--r--,"",d,17,root,root,3640,/usr/share/man/man8/cron.8.gz
rw-r--r--,"",d,18,root,root,36,/usr/share/man/man8/crond.8.gz
