# test/data/firewalld-0.2.12-5.fc18.noarch.filelist
rw-r--r--,"",cn,1,root,root,1247,/etc/dbus-1/system.d/FirewallD.conf
rw-r-----,"",cn,3,root,root,476,/etc/firewalld/firewalld.conf
rw-r--r--,"",cn,7,root,root,72,/etc/sysconfig/firewalld
rwxr-xr-x,"","",8,root,root,19080,/usr/bin/firewall-cmd
rwxr-xr-x,"","",9,root,root,16206,/usr/bin/firewall-offline-cmd
rw-r-----,"","",12,root,root,222,/usr/lib/firewalld/icmptypes/destination-unreachable.xml
rw-r-----,"","",13,root,root,192,/usr/lib/firewalld/icmptypes/echo-reply.xml
rw-r-----,"","",14,root,root,229,/usr/lib/firewalld/icmptypes/echo-request.xml
rw-r-----,"","",15,root,root,225,/usr/lib/firewalld/icmptypes/parameter-problem.xml
rw-r-----,"","",16,root,root,185,/usr/lib/firewalld/icmptypes/redirect.xml
rw-r-----,"","",17,root,root,227,/usr/lib/firewalld/icmptypes/router-advertisement.xml
rw-r-----,"","",18,root,root,241,/usr/lib/firewalld/icmptypes/router-solicitation.xml
rw-r-----,"","",19,root,root,248,/usr/lib/firewalld/icmptypes/source-quench.xml
rw-r-----,"","",20,root,root,253,/usr/lib/firewalld/icmptypes/time-exceeded.xml
rw-r-----,"","",22,root,root,374,/usr/lib/firewalld/services/amanda-client.xml
rw-r-----,"","",23,root,root,320,/usr/lib/firewalld/services/bacula-client.xml
rw-r-----,"","",24,root,root,346,/usr/lib/firewalld/services/bacula.xml
rw-r-----,"","",25,root,root,455,/usr/lib/firewalld/services/cluster-suite.xml
rw-r-----,"","",26,root,root,338,/usr/lib/firewalld/services/dhcpv6-client.xml
rw-r-----,"","",27,root,root,346,/usr/lib/firewalld/services/dns.xml
rw-r-----,"","",28,root,root,374,/usr/lib/firewalld/services/ftp.xml
rw-r-----,"","",29,root,root,353,/usr/lib/firewalld/services/http.xml
rw-r-----,"","",30,root,root,448,/usr/lib/firewalld/services/https.xml
rw-r-----,"","",31,root,root,372,/usr/lib/firewalld/services/imaps.xml
rw-r-----,"","",32,root,root,454,/usr/lib/firewalld/services/ipp-client.xml
rw-r-----,"","",33,root,root,427,/usr/lib/firewalld/services/ipp.xml
rw-r-----,"","",34,root,root,517,/usr/lib/firewalld/services/ipsec.xml
rw-r-----,"","",35,root,root,385,/usr/lib/firewalld/services/libvirt-tls.xml
rw-r-----,"","",36,root,root,389,/usr/lib/firewalld/services/libvirt.xml
rw-r-----,"","",37,root,root,424,/usr/lib/firewalld/services/mdns.xml
rw-r-----,"","",38,root,root,324,/usr/lib/firewalld/services/nfs.xml
rw-r-----,"","",39,root,root,389,/usr/lib/firewalld/services/ntp.xml
rw-r-----,"","",40,root,root,335,/usr/lib/firewalld/services/openvpn.xml
rw-r-----,"","",41,root,root,357,/usr/lib/firewalld/services/pop3s.xml
rw-r-----,"","",42,root,root,446,/usr/lib/firewalld/services/radius.xml
rw-r-----,"","",43,root,root,384,/usr/lib/firewalld/services/samba-client.xml
rw-r-----,"","",44,root,root,461,/usr/lib/firewalld/services/samba.xml
rw-r-----,"","",45,root,root,550,/usr/lib/firewalld/services/smtp.xml
rw-r-----,"","",46,root,root,463,/usr/lib/firewalld/services/ssh.xml
rw-r-----,"","",47,root,root,393,/usr/lib/firewalld/services/telnet.xml
rw-r-----,"","",48,root,root,301,/usr/lib/firewalld/services/tftp-client.xml
rw-r-----,"","",49,root,root,437,/usr/lib/firewalld/services/tftp.xml
rw-r-----,"","",50,root,root,475,/usr/lib/firewalld/services/vnc-server.xml
rw-r-----,"","",52,root,root,256,/usr/lib/firewalld/zones/block.xml
rw-r-----,"","",53,root,root,293,/usr/lib/firewalld/zones/dmz.xml
rw-r-----,"","",54,root,root,226,/usr/lib/firewalld/zones/drop.xml
rw-r-----,"","",55,root,root,319,/usr/lib/firewalld/zones/external.xml
rw-r-----,"","",56,root,root,400,/usr/lib/firewalld/zones/home.xml
rw-r-----,"","",57,root,root,415,/usr/lib/firewalld/zones/internal.xml
rw-r-----,"","",58,root,root,340,/usr/lib/firewalld/zones/public.xml
rw-r-----,"","",59,root,root,179,/usr/lib/firewalld/zones/trusted.xml
rw-r-----,"","",60,root,root,367,/usr/lib/firewalld/zones/work.xml
rw-r--r--,"","",62,root,root,0,/usr/lib/python2.7/site-packages/firewall/__init__.py
rw-r--r--,"","",63,root,root,140,/usr/lib/python2.7/site-packages/firewall/__init__.pyc
rw-r--r--,"","",63,root,root,140,/usr/lib/python2.7/site-packages/firewall/__init__.pyo
rw-r--r--,"","",65,root,root,30129,/usr/lib/python2.7/site-packages/firewall/client.py
rw-r--r--,"","",66,root,root,40735,/usr/lib/python2.7/site-packages/firewall/client.pyc
rw-r--r--,"","",66,root,root,40735,/usr/lib/python2.7/site-packages/firewall/client.pyo
rw-r--r--,"","",69,root,root,2748,/usr/lib/python2.7/site-packages/firewall/config/__init__.py
rw-r--r--,"","",70,root,root,2284,/usr/lib/python2.7/site-packages/firewall/config/__init__.pyc
rw-r--r--,"","",70,root,root,2284,/usr/lib/python2.7/site-packages/firewall/config/__init__.pyo
rw-r--r--,"","",72,root,root,1720,/usr/lib/python2.7/site-packages/firewall/config/dbus.py
rw-r--r--,"","",73,root,root,1010,/usr/lib/python2.7/site-packages/firewall/config/dbus.pyc
rw-r--r--,"","",73,root,root,1010,/usr/lib/python2.7/site-packages/firewall/config/dbus.pyo
rw-r--r--,"","",76,root,root,0,/usr/lib/python2.7/site-packages/firewall/core/__init__.py
rw-r--r--,"","",77,root,root,145,/usr/lib/python2.7/site-packages/firewall/core/__init__.pyc
rw-r--r--,"","",77,root,root,145,/usr/lib/python2.7/site-packages/firewall/core/__init__.pyo
rw-r--r--,"","",79,root,root,1724,/usr/lib/python2.7/site-packages/firewall/core/base.py
rw-r--r--,"","",80,root,root,841,/usr/lib/python2.7/site-packages/firewall/core/base.pyc
rw-r--r--,"","",80,root,root,841,/usr/lib/python2.7/site-packages/firewall/core/base.pyo
rw-r--r--,"","",82,root,root,2521,/usr/lib/python2.7/site-packages/firewall/core/ebtables.py
rw-r--r--,"","",83,root,root,2917,/usr/lib/python2.7/site-packages/firewall/core/ebtables.pyc
rw-r--r--,"","",83,root,root,2917,/usr/lib/python2.7/site-packages/firewall/core/ebtables.pyo
rw-r--r--,"","",85,root,root,20506,/usr/lib/python2.7/site-packages/firewall/core/fw.py
rw-r--r--,"","",86,root,root,16491,/usr/lib/python2.7/site-packages/firewall/core/fw.pyc
rw-r--r--,"","",86,root,root,16491,/usr/lib/python2.7/site-packages/firewall/core/fw.pyo
rw-r--r--,"","",88,root,root,18685,/usr/lib/python2.7/site-packages/firewall/core/fw_config.py
rw-r--r--,"","",89,root,root,16972,/usr/lib/python2.7/site-packages/firewall/core/fw_config.pyc
rw-r--r--,"","",89,root,root,16972,/usr/lib/python2.7/site-packages/firewall/core/fw_config.pyo
rw-r--r--,"","",91,root,root,7769,/usr/lib/python2.7/site-packages/firewall/core/fw_direct.py
rw-r--r--,"","",92,root,root,7084,/usr/lib/python2.7/site-packages/firewall/core/fw_direct.pyc
rw-r--r--,"","",92,root,root,7084,/usr/lib/python2.7/site-packages/firewall/core/fw_direct.pyo
rw-r--r--,"","",94,root,root,1534,/usr/lib/python2.7/site-packages/firewall/core/fw_icmptype.py
rw-r--r--,"","",95,root,root,2083,/usr/lib/python2.7/site-packages/firewall/core/fw_icmptype.pyc
rw-r--r--,"","",95,root,root,2083,/usr/lib/python2.7/site-packages/firewall/core/fw_icmptype.pyo
rw-r--r--,"","",97,root,root,1510,/usr/lib/python2.7/site-packages/firewall/core/fw_service.py
rw-r--r--,"","",98,root,root,2063,/usr/lib/python2.7/site-packages/firewall/core/fw_service.pyc
rw-r--r--,"","",98,root,root,2063,/usr/lib/python2.7/site-packages/firewall/core/fw_service.pyo
rw-r--r--,"","",100,root,root,30234,/usr/lib/python2.7/site-packages/firewall/core/fw_zone.py
rw-r--r--,"","",101,root,root,25334,/usr/lib/python2.7/site-packages/firewall/core/fw_zone.pyc
rw-r--r--,"","",101,root,root,25334,/usr/lib/python2.7/site-packages/firewall/core/fw_zone.pyo
rw-r--r--,"","",104,root,root,1586,/usr/lib/python2.7/site-packages/firewall/core/io/__init__.py
rw-r--r--,"","",105,root,root,917,/usr/lib/python2.7/site-packages/firewall/core/io/__init__.pyc
rw-r--r--,"","",105,root,root,917,/usr/lib/python2.7/site-packages/firewall/core/io/__init__.pyo
rw-r--r--,"","",107,root,root,6316,/usr/lib/python2.7/site-packages/firewall/core/io/firewalld_conf.py
rw-r--r--,"","",108,root,root,4799,/usr/lib/python2.7/site-packages/firewall/core/io/firewalld_conf.pyc
rw-r--r--,"","",108,root,root,4799,/usr/lib/python2.7/site-packages/firewall/core/io/firewalld_conf.pyo
rw-r--r--,"","",110,root,root,4994,/usr/lib/python2.7/site-packages/firewall/core/io/icmptype.py
rw-r--r--,"","",111,root,root,4554,/usr/lib/python2.7/site-packages/firewall/core/io/icmptype.pyc
rw-r--r--,"","",111,root,root,4554,/usr/lib/python2.7/site-packages/firewall/core/io/icmptype.pyo
rw-r--r--,"","",113,root,root,7820,/usr/lib/python2.7/site-packages/firewall/core/io/io_object.py
rw-r--r--,"","",114,root,root,9243,/usr/lib/python2.7/site-packages/firewall/core/io/io_object.pyc
rw-r--r--,"","",114,root,root,9243,/usr/lib/python2.7/site-packages/firewall/core/io/io_object.pyo
rw-r--r--,"","",116,root,root,6207,/usr/lib/python2.7/site-packages/firewall/core/io/service.py
rw-r--r--,"","",117,root,root,5272,/usr/lib/python2.7/site-packages/firewall/core/io/service.pyc
rw-r--r--,"","",117,root,root,5272,/usr/lib/python2.7/site-packages/firewall/core/io/service.pyo
rw-r--r--,"","",119,root,root,9471,/usr/lib/python2.7/site-packages/firewall/core/io/zone.py
rw-r--r--,"","",120,root,root,7225,/usr/lib/python2.7/site-packages/firewall/core/io/zone.pyc
rw-r--r--,"","",120,root,root,7225,/usr/lib/python2.7/site-packages/firewall/core/io/zone.pyo
rw-r--r--,"","",122,root,root,6514,/usr/lib/python2.7/site-packages/firewall/core/ipXtables.py
rw-r--r--,"","",123,root,root,6878,/usr/lib/python2.7/site-packages/firewall/core/ipXtables.pyc
rw-r--r--,"","",123,root,root,6878,/usr/lib/python2.7/site-packages/firewall/core/ipXtables.pyo
rw-r--r--,"","",125,root,root,29710,/usr/lib/python2.7/site-packages/firewall/core/logger.py
rw-r--r--,"","",126,root,root,28655,/usr/lib/python2.7/site-packages/firewall/core/logger.pyc
rw-r--r--,"","",126,root,root,28655,/usr/lib/python2.7/site-packages/firewall/core/logger.pyo
rw-r--r--,"","",128,root,root,3474,/usr/lib/python2.7/site-packages/firewall/core/modules.py
rw-r--r--,"","",129,root,root,3656,/usr/lib/python2.7/site-packages/firewall/core/modules.pyc
rw-r--r--,"","",129,root,root,3656,/usr/lib/python2.7/site-packages/firewall/core/modules.pyo
rw-r--r--,"","",131,root,root,1493,/usr/lib/python2.7/site-packages/firewall/core/prog.py
rw-r--r--,"","",132,root,root,1011,/usr/lib/python2.7/site-packages/firewall/core/prog.pyc
rw-r--r--,"","",132,root,root,1011,/usr/lib/python2.7/site-packages/firewall/core/prog.pyo
rw-r--r--,"","",134,root,root,2749,/usr/lib/python2.7/site-packages/firewall/core/watcher.py
rw-r--r--,"","",135,root,root,2767,/usr/lib/python2.7/site-packages/firewall/core/watcher.pyc
rw-r--r--,"","",135,root,root,2767,/usr/lib/python2.7/site-packages/firewall/core/watcher.pyo
rw-r--r--,"","",137,root,root,3786,/usr/lib/python2.7/site-packages/firewall/dbus_utils.py
rw-r--r--,"","",138,root,root,3981,/usr/lib/python2.7/site-packages/firewall/dbus_utils.pyc
rw-r--r--,"","",138,root,root,3981,/usr/lib/python2.7/site-packages/firewall/dbus_utils.pyo
rw-r--r--,"","",140,root,root,2885,/usr/lib/python2.7/site-packages/firewall/errors.py
rw-r--r--,"","",141,root,root,2972,/usr/lib/python2.7/site-packages/firewall/errors.pyc
rw-r--r--,"","",141,root,root,2972,/usr/lib/python2.7/site-packages/firewall/errors.pyo
rw-r--r--,"","",143,root,root,6295,/usr/lib/python2.7/site-packages/firewall/functions.py
rw-r--r--,"","",144,root,root,5916,/usr/lib/python2.7/site-packages/firewall/functions.pyc
rw-r--r--,"","",144,root,root,5916,/usr/lib/python2.7/site-packages/firewall/functions.pyo
rw-r--r--,"","",147,root,root,0,/usr/lib/python2.7/site-packages/firewall/server/__init__.py
rw-r--r--,"","",148,root,root,147,/usr/lib/python2.7/site-packages/firewall/server/__init__.pyc
rw-r--r--,"","",148,root,root,147,/usr/lib/python2.7/site-packages/firewall/server/__init__.pyo
rw-r--r--,"","",150,root,root,12439,/usr/lib/python2.7/site-packages/firewall/server/config.py
rw-r--r--,"","",151,root,root,11991,/usr/lib/python2.7/site-packages/firewall/server/config.pyc
rw-r--r--,"","",151,root,root,11991,/usr/lib/python2.7/site-packages/firewall/server/config.pyo
rw-r--r--,"","",153,root,root,7685,/usr/lib/python2.7/site-packages/firewall/server/config_icmptype.py
rw-r--r--,"","",154,root,root,7535,/usr/lib/python2.7/site-packages/firewall/server/config_icmptype.pyc
rw-r--r--,"","",154,root,root,7535,/usr/lib/python2.7/site-packages/firewall/server/config_icmptype.pyo
rw-r--r--,"","",156,root,root,7638,/usr/lib/python2.7/site-packages/firewall/server/config_service.py
rw-r--r--,"","",157,root,root,7480,/usr/lib/python2.7/site-packages/firewall/server/config_service.pyc
rw-r--r--,"","",157,root,root,7480,/usr/lib/python2.7/site-packages/firewall/server/config_service.pyo
rw-r--r--,"","",159,root,root,7496,/usr/lib/python2.7/site-packages/firewall/server/config_zone.py
rw-r--r--,"","",160,root,root,7343,/usr/lib/python2.7/site-packages/firewall/server/config_zone.pyc
rw-r--r--,"","",160,root,root,7343,/usr/lib/python2.7/site-packages/firewall/server/config_zone.pyo
rw-r--r--,"","",162,root,root,2095,/usr/lib/python2.7/site-packages/firewall/server/decorators.py
rw-r--r--,"","",163,root,root,1694,/usr/lib/python2.7/site-packages/firewall/server/decorators.pyc
rw-r--r--,"","",163,root,root,1694,/usr/lib/python2.7/site-packages/firewall/server/decorators.pyo
rw-r--r--,"","",165,root,root,38942,/usr/lib/python2.7/site-packages/firewall/server/firewalld.py
rw-r--r--,"","",166,root,root,34854,/usr/lib/python2.7/site-packages/firewall/server/firewalld.pyc
rw-r--r--,"","",166,root,root,34854,/usr/lib/python2.7/site-packages/firewall/server/firewalld.pyo
rw-r--r--,"","",168,root,root,2879,/usr/lib/python2.7/site-packages/firewall/server/server.py
rw-r--r--,"","",169,root,root,2042,/usr/lib/python2.7/site-packages/firewall/server/server.pyc
rw-r--r--,"","",169,root,root,2042,/usr/lib/python2.7/site-packages/firewall/server/server.pyo
rw-r--r--,"","",171,root,root,607,/usr/lib/systemd/system/firewalld.service
rwxr-xr-x,"","",172,root,root,3786,/usr/sbin/firewalld
rw-r--r--,"",d,174,root,root,18092,/usr/share/doc/firewalld-0.2.12/COPYING
rw-r--r--,"","",175,root,root,8818,/usr/share/locale/ar/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",176,root,root,10528,/usr/share/locale/as/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",177,root,root,12579,/usr/share/locale/bn_IN/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",178,root,root,7742,/usr/share/locale/ca/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",179,root,root,16766,/usr/share/locale/cs/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",180,root,root,7715,/usr/share/locale/da/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",181,root,root,18188,/usr/share/locale/de/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",182,root,root,7572,/usr/share/locale/en_GB/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",183,root,root,17757,/usr/share/locale/es/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",184,root,root,9640,/usr/share/locale/fi/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",185,root,root,18685,/usr/share/locale/fr/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",186,root,root,12104,/usr/share/locale/gu/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",187,root,root,11278,/usr/share/locale/hi/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",188,root,root,18370,/usr/share/locale/hu/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",189,root,root,17840,/usr/share/locale/it/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",190,root,root,15313,/usr/share/locale/ja/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",191,root,root,13739,/usr/share/locale/kn/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",192,root,root,8403,/usr/share/locale/ko/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",193,root,root,15226,/usr/share/locale/ml/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",194,root,root,11448,/usr/share/locale/mr/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",195,root,root,17557,/usr/share/locale/nl/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",196,root,root,13228,/usr/share/locale/or/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",197,root,root,18654,/usr/share/locale/pa/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",198,root,root,17992,/usr/share/locale/pl/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",199,root,root,8001,/usr/share/locale/pt/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",200,root,root,11023,/usr/share/locale/pt_BR/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",201,root,root,22427,/usr/share/locale/ru/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",202,root,root,7793,/usr/share/locale/sk/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",203,root,root,9894,/usr/share/locale/sr/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",204,root,root,7808,/usr/share/locale/sr@latin/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",205,root,root,17730,/usr/share/locale/sv/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",206,root,root,12890,/usr/share/locale/ta/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",207,root,root,11842,/usr/share/locale/te/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",208,root,root,2716,/usr/share/locale/tr/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",209,root,root,22852,/usr/share/locale/uk/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",210,root,root,16595,/usr/share/locale/zh_CN/LC_MESSAGES/firewalld.mo
rw-r--r--,"","",211,root,root,7571,/usr/share/locale/zh_TW/LC_MESSAGES/firewalld.mo
rw-r--r--,"",d,212,root,root,3680,/usr/share/man/man1/firewall-cmd.1.gz
rw-r--r--,"",d,213,root,root,2076,/usr/share/man/man1/firewalld.1.gz
rw-r--r--,"",d,214,root,root,783,/usr/share/man/man5/firewalld.conf.5.gz
rw-r--r--,"",d,215,root,root,886,/usr/share/man/man5/firewalld.icmptype.5.gz
rw-r--r--,"",d,216,root,root,1285,/usr/share/man/man5/firewalld.service.5.gz
rw-r--r--,"",d,217,root,root,1594,/usr/share/man/man5/firewalld.zone.5.gz
rw-r--r--,"",d,218,root,root,1687,/usr/share/man/man5/firewalld.zones.5.gz
rw-r--r--,"","",219,root,root,1428,/usr/share/polkit-1/actions/org.fedoraproject.FirewallD1.policy
